# pic_uploader.timer
[Unit]
Description=Run the pic_uploader script every 10 minutes

[Timer]
OnUnitActiveSec=30min
OnBootSec=30min

[Install]
WantedBy=timers.target
